# scheduler/platform/macos/com.cc.director.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cc.director</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>-m</string>
        <string>cc_director.service</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/cc_director/scheduler</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/path/to/cc_director/scheduler/logs/service_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/path/to/cc_director/scheduler/logs/service_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/path/to/cc_director/scheduler</string>
    </dict>
</dict>
</plist>
